# docs/systemd/antigravity-beat.service
# /etc/systemd/system/antigravity-beat.service
[Unit]
Description=Sanaa AI Celery Beat Scheduler
After=redis-server.service

[Service]
Type=exec
User=root
Group=antigravity
WorkingDirectory=/opt/antigravity/core
Environment="PATH=/opt/antigravity/venv/bin:/usr/bin"
EnvironmentFile=/opt/antigravity/.env
ExecStart=/opt/antigravity/venv/bin/celery -A tasks beat --loglevel=info
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
